# demo-proguy-me.service
[Unit]
Description=FastAPI site guard for demo.proguy.me
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=proguy
Group=proguy
WorkingDirectory=/home/proguy/demo.proguy.me
EnvironmentFile=/home/proguy/demo.proguy.me/.env
ExecStart=/home/proguy/demo.proguy.me/.venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8881
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/home/proguy/demo.proguy.me

[Install]
WantedBy=multi-user.target